# S73_HAD/lib/S73a2/app/lib/systemd/system/hobot_pms.service
[Unit]
Description=Hobot Power Management Service

[Service]
ExecStart=hb_pms

[Install]
WantedBy=multi-user.target